# systemd/parserUserBot-trigger.service
[Unit]
Description=parserUserBot — trigger worker (produces callback outbox entries)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=oneshot
User=kisskin
WorkingDirectory=/mnt/ssd/projects/parserUserBot
ExecStart=/usr/bin/python3 /mnt/ssd/projects/parserUserBot/trigger_worker.py --once
EnvironmentFile=/mnt/ssd/projects/parserUserBot/.env

StandardOutput=journal
StandardError=journal
