# systemd/parserUserBot-main.service
[Unit]
Description=parserUserBot — Telegram client daemon
After=network-online.target postgresql.service
Wants=network-online.target postgresql.service

[Service]
Type=simple
User=parseruserbot
WorkingDirectory=/opt/parserUserBot
EnvironmentFile=/opt/parserUserBot/.env
ExecStart=/opt/parserUserBot/.venv/bin/python /opt/parserUserBot/main.py
Restart=always
RestartSec=10
KillSignal=SIGTERM
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
