# systemd/parserUserBot-trigger.service
[Unit]
Description=parserUserBot — trigger worker (produces callback outbox entries)
After=network-online.target postgresql.service
Wants=network-online.target postgresql.service

[Service]
Type=oneshot
User=parseruserbot
WorkingDirectory=/opt/parserUserBot
EnvironmentFile=/opt/parserUserBot/.env
ExecStart=/opt/parserUserBot/.venv/bin/python /opt/parserUserBot/trigger_worker.py --once
StandardOutput=journal
StandardError=journal
