# launchd/com.kisskin.parseruserbot.trigger.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kisskin.parseruserbot.trigger</string>
    <key>WorkingDirectory</key>
    <string>/Users/your-user/parserUserBot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/your-user/parserUserBot/.venv/bin/python</string>
        <string>/Users/your-user/parserUserBot/trigger_worker.py</string>
        <string>--once</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>30</integer>
    <key>StandardOutPath</key>
    <string>/Users/your-user/parserUserBot/logs/trigger.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/your-user/parserUserBot/logs/trigger.err.log</string>
</dict>
</plist>
